# deploy/zul-health.timer
# Runs the Zul mainnet health check every 5 minutes.
[Unit]
Description=Run the Zul mainnet health check periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
